# provision/service/fastcgi-mono-server.conf
[program:fastcgi-mono-server]
command=fastcgi-mono-server4 /applications=/:/share/public/ /socket=tcp:127.0.0.1:9000
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true